# docker-file/config/supervisor/nginx-php.conf
[supervisord]
nodaemon=true
user = root
[program:php-fpm]
command=/usr/local/sbin/php-fpm -FR -c /usr/local/etc/php-fpm.conf
autostart = true
priority = 5
;command=/usr/sbin/php70-fpm -c /etc/php-fpm.d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart = true
priority = 10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0